# thermal/temperature.service
[Unit]
Description=A job that stores the temperature
[Service]
Type=oneshot
ExecStart=/opt/edgecaffe/thermal/record.sh